# supervisor/conf.d/consumer.conf
[program:consumer]
command = /usr/bin/env python ./consumer.py
stdout_logfile = /tmp/consumer.stdout.log
stderr_logfile = /tmp/consumer.stderr.log
stopasgroup = true
