# supervisor/conf.d/intake-server.conf
[program:intake-server]
command = /usr/bin/env intake-server ./catalog.yml
stdout_logfile = /tmp/intake-server.stdout.log
stderr_logfile = /tmp/intake-server.stderr.log
stopasgroup = true
