# supervisor/conf.d/mongobox.conf
[program:consumer]
command = /usr/bin/env python ./run-mongobox.py
stdout_logfile = /tmp/run-mongobox.stdout.log
stderr_logfile = /tmp/run-mongobox.stderr.log
stopasgroup = true
